# k11-boxing.service
[Unit]
Description=K11 Boxing - Servidor Node.js
Documentation=https://github.com/K11-Boxing
After=network.target network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/k11-boxing
ExecStart=/usr/local/bin/node --max-old-space-size=256 --v8-pool-size=4 /home/pi/k11-boxing/server.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=DISPLAY=:0
Nice=-5

# Seguridad minima (no root)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
